# scripts/datenschleuder_agent/datenschleuder-agent.service
[Unit]
Description=Datenschleuder Agent for Datenschleuder
After=network.target

[Service]
Type=simple
User=datenschleuder-agent
Group=docker
WorkingDirectory=/opt/datenschleuder-agent
ExecStart=/opt/datenschleuder-agent/venv/bin/python agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
EnvironmentFile=/opt/datenschleuder-agent/.env

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
